# infra/systemd/kalshi-bot-compose.service
[Unit]
Description=Kalshi Bot Docker Compose Stack (web and daemon)
After=network-online.target docker.service
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/workspace/kalshi_bot
ExecStart=/bin/bash -lc './infra/scripts/start-stack.sh systemd_boot'
ExecStop=/bin/bash -lc '/usr/bin/docker compose -f infra/docker-compose.yml down'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
